# ytinterface.service
[Unit]
Description=ytinterface - YouTube frontend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ytinterface
Group=ytinterface
WorkingDirectory=/opt/ytinterface
Environment=FLASK_DEBUG=0
Environment=PATH=/opt/ytinterface/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/ytinterface/.venv/bin/python app.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
